# deploy/acme/hyproxy-acme.service
# Issue/renew the wildcard Let's Encrypt cert via DNS-01 and install it into the
# data-plane cert paths. Oneshot, driven by hyproxy-acme.timer.
#
# Install: copy to /etc/systemd/system/, adjust paths, then
#   systemctl daemon-reload && systemctl enable --now hyproxy-acme.timer
# Provider creds + HYPROXY_DOMAIN/ACME_EMAIL/LEGO_DNS_PROVIDER live in
# /etc/hyproxy/acme.env (root-owned, 0600).
[Unit]
Description=hyproxy ACME wildcard cert issuance/renewal (lego DNS-01)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/hyproxy/acme.env
# Adjust to the deployed repo path:
ExecStart=/opt/hyproxy/deploy/acme/obtain-cert.sh
