# deploy/acme/hyproxy-acme.timer
# Daily cert-renewal check (lego renews only when <30 days remain, so a daily
# run with jitter is safe and self-throttling).
[Unit]
Description=Daily hyproxy Let's Encrypt renewal check

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
